# mission-control-sync.service
[Unit]
Description=Mission Control Cron State Sync
After=network.target

[Service]
Type=simple
ExecStart=/home/linuxbrew/.linuxbrew/bin/node /home/spindux/mission-control/sync-cron-state.js
Restart=always
RestartSec=10
User=spindux

# Environment
Environment=HOME=/home/spindux
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/spindux/.local/bin
Environment=NODE_ENV=production

# Logging
StandardOutput=append:/home/spindux/mission-control/sync.log
StandardError=append:/home/spindux/mission-control/sync-error.log

[Install]
WantedBy=default.target